# etc/Telegram.conf
Description "Telegram UpStart Script"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

setuid cli
exec /bin/sh cli/launch.sh
